# systemd/efistub-update.path
[Unit]
Description=Update efistub boot config when kernel files are changed

[Path]
PathChanged=/boot/initramfs-linux-fallback.img

[Install]
WantedBy=multi-user.target
WantedBy=system-update.target
