# systemd/efistub-update.service
[Unit]
Description=Update efistub boot config when kernel files are updated

[Service]
Type=oneshot
ExecStart=/usr/bin/efistub bootctl update
